# opt/biotech-dashboard/etc/systemd/system/biotech-collector.service
# /etc/systemd/system/biotech-collector.service
[Unit]
Description=Biotech Data Collection Service

[Service]
Type=oneshot
User=biotech
WorkingDirectory=/opt/biotech-dashboard/scripts
Environment="PATH=/opt/biotech-dashboard/venv/bin"
EnvironmentFile=/opt/biotech-dashboard/.env
ExecStart=/opt/biotech-dashboard/venv/bin/python polygon_fetch.py